# brain/deploy/systemd/triage-brain-app.service
[Unit]
Description=TRIAGE Brain Pi app server
After=network.target triage-brain-camera.service
Wants=triage-brain-camera.service

[Service]
Type=simple
User=__TRIAGE_USER__
WorkingDirectory=__TRIAGE_ROOT__
EnvironmentFile=__TRIAGE_ROOT__/.env
ExecStart=/usr/bin/node __TRIAGE_ROOT__/brain/server/index.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
