# brain/deploy/systemd/triage-brain-camera.service
[Unit]
Description=TRIAGE Brain Pi camera service
After=network.target

[Service]
Type=simple
User=__TRIAGE_USER__
WorkingDirectory=__TRIAGE_ROOT__
EnvironmentFile=__TRIAGE_ROOT__/.env
ExecStart=/usr/bin/python3 __TRIAGE_ROOT__/brain/camera/camera.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
